# Library/LaunchAgents/org.freedesktop.xdg-environment-vars.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.freedesktop.xdg-environment-vars</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>
    # Track existing environment variables
    EXISTING_ENV=$(awk 'BEGIN { for (name in ENVIRON) print name }' | sort)

    # Load the XDG environment variables
    set -a
    source "$HOME/.config/xdg.env"
    set +a

    # Add all xdg.env environment vars to launchctl
    for VAR in $(awk 'BEGIN { for (name in ENVIRON) print name }' | sort | comm -13 &lt;(printf '%s\n' "$EXISTING_ENV") -); do
      launchctl setenv "$VAR" "${!VAR}"
    done

    # Ensure XDG directories all exist
    mkdir -p "$XDG_CACHE_HOME" "$XDG_CONFIG_HOME" "$XDG_DATA_HOME" "$XDG_RUNTIME_DIR" "$XDG_STATE_HOME"
    </string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
